# packaging/linux/auto_pull.service
[Unit]
Description=auto_pull git watcher daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/auto_pull /etc/auto_pull/config_auto_pull.json
Restart=always
RestartSec=5
WorkingDirectory=/etc/auto_pull
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
